# system/etc/rc.local
#!/bin/bash

#rm /var/lib/dhcp/dhclient.leases

dhclient -r
dhclient
systemctl restart systemd-networkd

echo 1 > /tmp/rc.local

cd /opt/CloudLocalNetAgent
/opt/CloudLocalNetAgent/dev_run.sh &
#/opt/CloudLocalNetAgent/dev_run.sh >> /opt/agentlog.txt 2>&1 &

exit 0